# files/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv sh
exec php-fpm7 --nodaemonize --force-stderr --pid /var/run/php-fpm.pid
